# fan_control.service
[Unit]
Description=Temperature controlled fan

[Service]
ExecStart=/usr/local/bin/fan_control

[Install]
WantedBy=multi-user.target
